# deploy/celery.service
[Unit]
Description=BuyruqSportedu Celery Worker
After=network.target redis.service postgresql.service
Requires=redis.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/buyruqsportedu/backend
EnvironmentFile=/var/www/buyruqsportedu/backend/.env
ExecStart=/var/www/buyruqsportedu/backend/venv/bin/celery \
    -A config worker \
    --loglevel=info \
    --concurrency=4 \
    --logfile=/var/log/buyruqsportedu/celery.log \
    --pidfile=/var/run/buyruqsportedu/celery.pid
ExecStop=/bin/kill -s TERM $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
